# supervisor/kairos.conf
[program:kairos]
command=/kairos-run.sh
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
